# snortalog.cron
#!/bin/sh
cat /var/log/snort/* /var/log/snort/*/* 2>&1 | /usr/bin/snortalog.pl -report 
