# snortalog.cron
#!/bin/sh
if [ -d /var/log/snort ]; then
	cat /var/log/snort/* /var/log/snort/*/* 2>&1 | /usr/bin/snortalog.pl -report 
fi
